# dist/relay-server.service
[Unit]
Description=selective-relay-server
After=network.target

[Service]
Type=simple
User=selective-relay
WorkingDirectory=/home/selective-relay/selective-relay
ExecStart=/home/selective-relay/selective-relay/bin/server
ExecReload=/bin/kill -SIGUSR1 $MAINPID
TimeoutSec=15
Restart=always

[Install]
WantedBy=multi-user.target
